# cronjobs/crontab

SHELL=/bin/bash

# run crontab every hours
0 * * * *   /usr/bin/docker ps | grep aptly

# “At 04:05 on Tuesday.”
5 4 * * 2 /cronjobs/aptly-update-mirrors.sh

