# cronjobs/crontab

SHELL=/bin/bash

# run crontab every hours
3 * * * *   /cronjobs/aptly-graph.sh

# “At 04:05 on Tuesday.”
5 4 * * 2 /cronjobs/aptly-update-mirrors.sh

